# backend/avaliador-api.service
[Unit]
Description=Avaliador Prompt API Service
After=network.target

[Service]
User=appuser
Group=appuser
WorkingDirectory=/var/www/avaliadorprompt/backend
ExecStart=/var/www/avaliadorprompt/venv/bin/python run_server.py
Restart=on-failure
RestartSec=10
Environment="PATH=/var/www/avaliadorprompt/venv/bin"
EnvironmentFile=/var/www/avaliadorprompt/backend/.env.production

[Install]
WantedBy=multi-user.target
